# res/cron/dns_server.crontab
# In general this needs to be manually updated on the nameservers whenever a
# change is made here.
## GET UPDATES AND LOAD THEM IF THEY ARE GOOD ##
10,40 * * * * svn update /etc/bind > /dev/null
15,45 * * * * /opt/csw/sbin/rndc reload > /dev/null
## Check the status of DNS before everyone goes home (m-f)
46 16 * * 1-5 /opt/csw/sbin/rndc status
05,35 * * * * /etc/bind/cron/update_dns > /dev/null
